# templates/cache/restic-cache-cleanup.timer
[Timer]
OnCalendar={{ restic.cache_cleanup.schedule }}

[Install]
WantedBy=timers.target
