# templates/cache/restic-cache-cleanup.timer
[Timer]
OnCalendar=*-*-* 06:21:00

[Install]
WantedBy=timers.target
