# templates/cache/restic-cache-cleanup.service
[Service]
Type=oneshot
ExecStart=/usr/local/bin/restic-cache-cleanup
