# templates/systemd/restic-backup.service
[Service]
Type=oneshot
ExecStart={{ restic.prefix.bin }}/restic-backup-{{ item.key }}
