# templates/systemd/restic-backup.timer
[Timer]
{{ restic_config.timer|default('', true) }}

[Install]
WantedBy=timers.target
